# infra/vm-discovery-2026-05-05/sync-lustyle-blogs.cron
#!/bin/sh
su lu -c 'cd /home/lu/blogs/source/_posts && /usr/bin/git pull origin master'
